# scripts/jamhq-automation-worker.service
[Unit]
Description=Jam HQ automation worker
After=network.target jam-hq.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/alpegames/jam-hq
ExecStart=/usr/bin/npx tsx scripts/worker.mjs
Restart=always
RestartSec=5
EnvironmentFile=/opt/alpegames/jam-hq/.env.worker

[Install]
WantedBy=multi-user.target
